# backend/systemd/celery-beat.service
[Unit]
Description=Celery Beat Scheduler for AlgoItny
Documentation=https://docs.celeryproject.org/
After=network.target redis.service

[Service]
Type=simple
User=algoitny
Group=www-data
WorkingDirectory=/home/algoitny/apps/algoitny/backend
Environment="PATH=/home/algoitny/apps/algoitny/backend/venv/bin"
EnvironmentFile=/home/algoitny/apps/algoitny/backend/.env

# Celery beat execution
ExecStart=/home/algoitny/apps/algoitny/backend/venv/bin/celery -A config beat \
    --loglevel=info \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler \
    --logfile=/var/log/celery/beat.log \
    --pidfile=/var/run/celery/beat.pid

# Process management
Restart=always
RestartSec=10
TimeoutStopSec=30

# Security settings
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
